# bot.conf
[program:rosatom_bot]
directory=/home/ubuntu/rosatom_bot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/rosatom_bot/logfile_err.log
stdout_logfile=/home/ubuntu/rosatom_bot/logfile.log